# packages/debian/http-remote-pc-control.service
[Unit]
Description=HTTP Remote PC Control
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/bin
ExecStart=/usr/bin/http-remote-pc-control
Restart=always

[Install]
WantedBy=multi-user.target
