# packages/http-remote-pc-control.service
[Unit]
Description=HTTP Remote PC Control
After=graphical-session.target network-online.target
Wants=network-online.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStartPre=/usr/bin/http-remote-pc-control --generate --if-missing
ExecStart=/usr/bin/http-remote-pc-control
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target