# packages/http-remote-pc-control.service
[Unit]
Description=HTTP Remote PC Control
After=graphical-session.target

[Service]
Type=simple
ExecStartPre=/usr/share/http-remote-pc-control/gen-certs
ExecStart=/usr/bin/http-remote-pc-control --port=5000 --cert-dir=%h/.local/share/http-remote-pc-control/certs
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
